# ops/vultr/systemd/cryptosquid-dashboard.service
[Unit]
Description=Crypto Squid Dashboard and Control API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cryptosquid
WorkingDirectory=/opt/cryptosquid
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/cryptosquid/.venv/bin/python src/dashboard.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
